# config/supervisor/apps.conf
[program:nginx-app]
command = /usr/sbin/nginx
